# .config/bash/profile.d/aliases.sh
#!/bin/bash

alias ec="emacsclient"
alias ff="open -a Firefox"
alias ll="ls -l"
